# systemd/claude-remote.service
# Long-lived Claude Code Remote Control session for the experiments workflow.
# Installed to ~/.config/systemd/user/claude-remote.service by bin/bootstrap.sh.
#
# Restart=always means that if the process dies or the network drops >10 min and the
# session times out, systemd brings it straight back. Because --name is fixed, the
# session always reappears under the same title in the phone's session list.
[Unit]
Description=Claude Code Remote Control (brot-os)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@@ROOT@@
ExecStart=@@CLAUDE@@ remote-control --name "brot-os"
Restart=always
RestartSec=15
# nvm-installed node must be on PATH for the session's tools.
Environment=PATH=@@NODEDIR@@:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
